# vagrant/puppet/modules/mongodb/templates/redhat_mongos-init.conf.erb
#!/bin/bash
# !!! generated by puppet !!!

# mongos - Startup script for mongos

# chkconfig: 35 85 15
# description: Mongo is a scalable, document-oriented database.
# processname: mongos
# config: /etc/mongos_<%= mongos_instance %>.conf

. /etc/rc.d/init.d/functions

# things from mongos.conf get there by mongos reading it


# NOTE: if you change any OPTIONS here, you get what you pay for:
# this script assumes all options are in the config file.
CONFIGFILE="/etc/mongos_<%= mongos_instance %>.conf"
OPTIONS=" -f $CONFIGFILE"

# FIXME: 1.9.x has a --shutdown flag that parses the config file and
# shuts down the correct running pid, but that's unavailable in 1.8
# for now.  This can go away when this script stops supporting 1.8.
PIDFILE=`awk -F= '/^pidfilepath=/{print $2}' "$CONFIGFILE"`
DBPATH=`dirname $PIDFILE`
mongos=/usr/bin/mongos

MONGO_USER=<%= scope.lookupvar('mongodb::params::run_as_user') %>
MONGO_GROUP=<%= scope.lookupvar('mongodb::params::run_as_group') %>

start()
{
  echo -n $"Starting mongos_<%= mongos_instance %>: "
  mkdir -p $DBPATH || exit 5
  chown -R ${MONGO_USER}:${MONGO_GROUP} $DBPATH
  daemon --user $MONGO_USER $mongos $OPTIONS
  RETVAL=$?
  echo
}

stop()
{
  echo -n $"Stopping mongos_<%= mongos_instance %>: "
  killproc -p "$PIDFILE" -d 30 $mongos -TERM
  RETVAL=$?
  rm "$PIDFILE"
  echo
}

restart () {
        stop
        start
}

ulimit -n 12000
RETVAL=0

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  status)
    status -p $PIDFILE
    RETVAL=$?
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload}"
    RETVAL=1
esac

exit $RETVAL