# docker.conf
[program:docker]
command=/usr/local/bin/docker -dD -s aufs
autostart=true
autorestart=true
stdout_logfile=/var/log/docker.log
redirect_stderr=true
numprocs=1
